# agentdvr/rootfs/etc/services.d/agentdvr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: AgentDVR
# Runs the AgentDVR Controller
# ==============================================================================




bashio::log.info 'Starting the AgentDVR...'
# Run agentdvr
exec ~/.dotnet/dotnet /agent/Agent.dll