# example/rootfs/etc/services.d/example1/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Home Assistant Add-on: Example
# Runs example1 script
# ==============================================================================


bashio::log.info "Starting Example1..."

exec /usr/bin/example1.sh
