# example/rootfs/etc/s6-overlay/s6-rc.d/example2/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs example2 script
# ==============================================================================

bashio::log.info "Starting Example2..."

exec /usr/bin/example2.sh
